# crontab
# This file runs the WeatherStation App every minute and writes results to a log file for debugging
*/1 * * * * python /root/WeatherStation.py > /tmp/test.log 2>&1
# This is all for now
